# systemd_pharma_news.service
[Unit]
Description=Pharmaceutical News Agent Daily Pipeline
After=network.target

[Service]
Type=oneshot
User=develop
Group=develop
WorkingDirectory=/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01
Environment="PATH=/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/config/.env
ExecStart=/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/venv/bin/python /home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/src/run_pipeline.py
StandardOutput=append:/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/logs/systemd_output.log
StandardError=append:/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/logs/systemd_error.log

[Install]
WantedBy=multi-user.target
