# systemd_pharma_admin_api.service
[Unit]
Description=Pharma News Admin API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=develop
Group=develop
WorkingDirectory=/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01
Environment="PATH=/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/config/.env
ExecStart=/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/venv/bin/python -m uvicorn src.admin_api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/logs/admin_api_output.log
StandardError=append:/home/develop/Daewoong_Pharmaceutical_News_Scrapper_PRODUCT01/logs/admin_api_error.log

[Install]
WantedBy=multi-user.target
